# ops/systemd/dailymacros-db-backup.timer
[Unit]
Description=Run DailyMacros logical Postgres backup before the daily EBS snapshot

[Timer]
OnCalendar=*-*-* 02:35:00
Persistent=true

[Install]
WantedBy=timers.target
